# cookbooks/openruko/templates/default/upstart-openruko-httprouting.conf.erb
start on starting openruko
stop on stopping openruko
respawn

setuid <%= node['user'] %>
setgid <%= node['group'] %>

env PG_USER=<%= node['user'] %>
env PG_PASSWORD=<%= node['postgresql']['password'][node['user']] %>
env HTTPROUTING_PORT=8080
env HTTPROUTING_TLS=true

chdir <%= node['openruko']['home'] %>/httprouting

script
  ./httprouting/bin/httprouting >> /var/log/openruko/httprouting.log 2>&1
end script
