# cookbooks/openruko/templates/default/upstart-openruko-httprouting.conf.erb
start on starting openruko-apiserver
stop on stopping openruko-apiserver
respawn

setuid root
setgid root

env PG_USER=<%= node['user'] %>
env PG_PASSWORD=<%= node['postgresql']['password'][node['user']] %>
env HTTPROUTING_PORT=80
env HTTPROUTING_TLS=true

chdir <%= node['openruko']['home'] %>/httprouting

script
  ./httprouting/bin/httprouting >> /var/log/openruko/httprouting.log 2>&1
end script
