# cookbooks/openruko/templates/default/upstart-openruko-dynohost.conf.erb
start on starting openruko
stop on stopping openruko
respawn

env RUKORUN_PATH="<%= node['openruko']['home'] %>/rukorun"
env CODONHOOKS_PATH="<%= node['openruko']['home'] %>/codonhooks"
env APISERVER_KEY="<%= node['openruko']['apiserver_key'] %>"

chdir <%= node['openruko']['home'] %>/dynohost

script
  ./dynohost/bin/dynohost >> /var/log/openruko/dynohost.log 2>&1
end script
