# cookbooks/openruko/templates/default/upstart-openruko-gitmouth.conf.erb
start on starting openruko
stop on stopping openruko
respawn

setuid <%= node['user'] %>
setgid <%= node['group'] %>

env APISERVER_KEY="<%= node['openruko']['apiserver_key'] %>"

chdir <%= node['openruko']['home'] %>/gitmouth

script
  ./bin/gitmouth >> /var/log/openruko/gitmouth.log 2>&1
end script
